# .config/systemd/user/battery-monitor.timer
[Unit]
Description=Run battery monitor every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
